# files/system/etc/profile.d/00-aliases.sh
#!/usr/bin/bash

alias ls="eza"
alias ll="eza -al"
alias l.="eza -d .*"