# files/system/etc/profile.d/zz-aliases.sh
#!/usr/bin/bash

alias ls="eza 2>/dev/null"
alias ll="eza -al 2>/dev/null"
alias l.="eza -d .* 2>/dev/null"
